# system_files/etc/skel/.config/autostart/tailscale-firstboot.desktop
[Desktop Entry]
Type=Application
Name=Tailscale First Boot
Exec=ptyxis -- /usr/bin/tailscale-firstboot.sh
X-GNOME-Autostart-enabled=true