# system_files/etc/skel/.config/autostart/tailscale-firstboot.desktop
[Desktop Entry]
Type=Application
Name=Tailscale First Boot
Exec=bash -c 'if [ ! -f /var/lib/tailscale/.firstboot ]; then ptyxis -- /usr/bin/tailscale-firstboot.sh; fi'
X-GNOME-Autostart-enabled=true